# deploy/rpmmanager.service
[Unit]
Description=RPM Manager - Graphical RPM Package Management Tool
After=network.target

[Service]
Type=simple
User=rpmmanager
Group=rpmmanager
ExecStart=/usr/local/bin/rpmmanager serve --config /etc/rpmmanager/config.yaml
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/rpmmanager /var/log/rpmmanager
PrivateTmp=true

[Install]
WantedBy=multi-user.target
